# deploy/trello.service.example
[Unit]
Description=Trello custom Next.js server
After=network.target

[Service]
Type=simple
User=__APP_USER__
WorkingDirectory=__APP_DIR__
Environment=NODE_ENV=production
Environment=PORT=3000
EnvironmentFile=__APP_DIR__/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
KillSignal=SIGINT
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
